# system/can2mqtt_hbay.service
[Unit]
Description=Hawkes Bay CAN to MQTT Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Wait until can0 is up
ExecStartPre=/bin/sh -c 'until ip link show can0 | grep -q "UP"; do sleep 1; done'
ExecStart=/usr/bin/python3 /home/pi/hbay/can2mqtt_hbay.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
